# src/zonst/qipai-sports/api/configapisrv/docs/srvapisrv.ini
[program:srvapisrv]
directory=/home/tonnn/srvapisrv/
command=/home/tonnn/srvapisrv/srvapisrv -config=/home/tonnn/srvapisrv/qcloud.toml
stdout_logfile=/data/log/srvapisrv/srvapisrv.log
stdout_logfile_backups=50
redirect_stderr=true
autostart=true
autorestart=true
